# etc/systemd/system/flask.service
[Unit]
Description=Flask-project
[Service]
User=jenkins
ExecStart=/var/lib/jenkins/workspace/flask_project/Script/installation.sh

[Install]
WantedBy=multi-user.target